# lab5/rootfs/etc/inittab
# Essential filesystem setup
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mount -t sysfs sysfs /sys
::sysinit:/bin/mount -t tmpfs tmpfs /run
::sysinit:/bin/mount -t tmpfs tmpfs /tmp
::sysinit:/bin/mount -o remount,rw /
::sysinit:/bin/mkdir -p /dev/pts /dev/shm
::sysinit:/bin/hostname -F /etc/hostname
::sysinit:/bin/mount -a
::sysinit:/bin/mount -t devpts devpts /dev/pts
::sysinit:/sbin/swapon -a
null::sysinit:/bin/ln -sf /proc/self/fd /dev/fd
null::sysinit:/bin/ln -sf /proc/self/fd/0 /dev/stdin
null::sysinit:/bin/ln -sf /proc/self/fd/1 /dev/stdout
null::sysinit:/bin/ln -sf /proc/self/fd/2 /dev/stderr

# now run any rc scripts
::wait:/etc/init.d/rcS

# Put a getty on the serial port
tty1::askfirst:/sbin/getty 38400 tty1 # Virtual Console 1 (e.g., HDMI)
tty2::askfirst:/sbin/getty 38400 tty2 # Virtual Console 2
::once:/sbin/getty -L ttyS0 115200 vt100 # BeagleBone Serial/UART0
::once:/sbin/getty -L ttyS4 115200 vt100 # BeagleBone Serial/UART4
::once:/sbin/getty -L ttyAMA0 115200 vt100 # QEMU Serial

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
::shutdown:/etc/init.d/rcK
::shutdown:/sbin/swapoff -a
::shutdown:/bin/sync
::shutdown:/bin/umount -a -r
